# layers/sway/sway.service
[Unit]
Description=Sway Wayland Compositor
After=dbus-session.service
Requires=dbus-session.service

[Service]
Type=simple
ExecStart=%h/.local/bin/sway-wrapper
Environment=XDG_RUNTIME_DIR=/tmp
Environment=WAYLAND_DISPLAY=wayland-0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/tmp/dbus-session
Environment=EGL_LOG_LEVEL=fatal
Environment=WLR_LIBINPUT_NO_DEVICES=1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
