# layers/google-chrome-sway/google-chrome.service
[Unit]
Description=Google Chrome Browser
After=sway.service
Requires=sway.service

[Service]
Type=simple
ExecStart=%h/.local/bin/chrome-wrapper
Environment=XDG_RUNTIME_DIR=/tmp
Environment=WAYLAND_DISPLAY=wayland-0
Environment=CHROME_FLAGS=--ozone-platform=wayland --enable-features=UseOzonePlatform
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
